# initrd/etc/init/init-system-dbus.conf
# init-system-dbus - tell init process to connect to system dbus server

start on started messagebus

task
exec /bin/kill -USR1 1
